# quadlet/jellyfin-local-trailers.container
[Unit]
Description = Container service for jellyfin-local-trailers

[Container]
ContainerName = %p
Image = jellyfin-local-trailers.build